# aria2.service
[Unit]
Description=Aria2 Service
After=network.target

[Service]
User=admin
ExecStart=/usr/bin/aria2c

[Install]
WantedBy=default.target
